# crontab
0 * * * * root bash /script.sh
@reboot root bash /script.sh
