# crontab
0 * * * * root bash /beem.sh
@reboot root bash /beem.sh
