# crontab
0 * * * * root bash /angers.sh
@reboot root bash /angers.sh
